# src/nginx/s6-rc.d/ticket-reload/run
#!/bin/sh
# 24h loop: rotate the TLS session-ticket key, then reload nginx IFF the config
# still tests clean. Lifted from the old bootstrap.sh tail; now a supervised
# longrun instead of a backgrounded `&`.
exec 2>&1
TICKET_KEY=/run/nginx/ticket.key
echo "[NGINX] 24h reload + TLS ticket-key rotation loop started (s6-supervised)."
while true; do
    sleep 1d
    [ -d /run/nginx ] && openssl rand 80 > "${TICKET_KEY}.tmp" 2>/dev/null \
        && chmod 0600 "${TICKET_KEY}.tmp" && mv -f "${TICKET_KEY}.tmp" "${TICKET_KEY}" || true
    if /usr/sbin/nginx -t 2>/dev/null; then
        /usr/sbin/nginx -s reload
    else
        echo "[NGINX] ---> 24h reload SKIPPED: nginx -t failed; keeping running config."
        /usr/sbin/nginx -t 2>&1 | sed 's/^/[NGINX] ---> /'
    fi
done
