# src/nginx/s6-rc.d/nginx/run
#!/bin/sh
# nginx longrun (web edge). Foreground; s6 supervises. Graceful stop via SIGQUIT
# configured in ./down-signal to match the image STOPSIGNAL. Applies the MALLOC
# LD_PRELOAD chosen by init-bootstrap.
exec 2>&1
[ -s /run/nginx/ld_preload ] && export LD_PRELOAD="$(cat /run/nginx/ld_preload)"
exec /usr/sbin/nginx -g "daemon off;"
