# src/roundcube/s6-rc.d/init-bootstrap/up
/etc/s6-overlay/scripts/init-bootstrap.sh
