# nordvpn2wireguard/rootfs/etc/services.d/nordvpn/run
#!/bin/bash

if [[ ! -d /run/nordvpn ]]; then
  mkdir -m 0770 /run/nordvpn
fi

exec s6-notifyoncheck -d /usr/sbin/nordvpnd > /dev/null
